# machine-learning/rootfs/etc/s6-overlay/s6-rc.d/longrun-server/run
#!/command/with-contenv ash
# shellcheck shell=bash

cd /app/

exec s6-setuidgid paperless-ml  \
    /usr/local/bin/granian \
        --host 0.0.0.0 \
        --port 8000 \
        --interface asgi \
        --no-ws \
        main:app
